# profile.d/manjaro.sh
#!/bin/sh

export MANJARO_LIB_MSG=/usr/lib/manjaro/util-msg.sh
export MANJARO_LIB_ALPM=/usr/lib/manjaro/util-alpm.sh
