# profile.d/manjaro.sh
#!/bin/sh

# MANJARO_MSG_LIB=/usr/lib/manjaro/util-msg.sh

export MANJARO_MSG_LIB=/usr/lib/manjaro/util-msg.sh
