# profile.d/manjaro.sh
#!/bin/sh

export MANJARO_MSG_LIB=/usr/lib/manjaro/util-msg.sh
export MANJARO_ALPM_LIB=/usr/lib/manjaro/util-alpm.sh
